# home-task-scheduler/rootfs/etc/services.d/home-task-scheduler/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

declare message
bashio::log.info "${message:="Starting Home Task Scheduler server..."}"

## Run your program
exec /usr/bin/run.sh
